# files/etc/systemd/system/margine-qemu-validation-inhibit.service
[Unit]
Description=Margine QEMU validation idle/suspend inhibitor
ConditionPathExists=/etc/margine/qemu-validation-inhibit.enabled

[Service]
Type=simple
ExecStart=/usr/bin/systemd-inhibit --what=sleep:idle --who=Margine-QEMU-validation "--why=QEMU validation in progress" --mode=block /usr/bin/sleep infinity
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target
