# files/home/.config/systemd/user/swaync.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/swaync --skip-system-css
