# files/etc/NetworkManager/dispatcher.d/90-margine-waybar-network
#!/usr/bin/env bash
set -euo pipefail

if [[ -x /usr/local/bin/margine-waybar-signal ]]; then
  /usr/local/bin/margine-waybar-signal 10
fi
