# files/etc/udev/rules.d/90-margine-waybar-power.rules
SUBSYSTEM=="power_supply", RUN+="/usr/local/bin/margine-waybar-signal 11"
