# files/home/.config/systemd/user/margine-maintenance-check.timer
[Unit]
Description=Periodic Margine update and advisory check

[Timer]
OnBootSec=10m
OnUnitActiveSec=6h
Persistent=true
Unit=margine-maintenance-check.service

[Install]
WantedBy=timers.target
