# files/home/.config/systemd/user/margine-notification-action-router.service
[Unit]
Description=Margine notification action router
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/margine-notification-action-router
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
